# docs/systemd/mesh-source-pacman.service
[Unit]
Description=Mesh source ingest: pacman log

[Service]
Type=oneshot
ExecStart=%h/.local/bin/meshctl event ingest-pacman --pacman-log /var/log/pacman.log
