# docs/systemd/mesh-source-cloudflared.timer
[Unit]
Description=Mesh source ingest timer: cloudflared journal

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
Persistent=true
Unit=mesh-source-cloudflared.service

[Install]
WantedBy=timers.target
